# files/snortcleanup.cron
# snort binary alert Log Cleanup
0 2   * * *     root  /usr/local/sbin/snortcleanup.sh
